# com.cortadogroup.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cortadogroup.watchdog</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/david/work/website_watchdog/venv/bin/python3</string>
        <string>/Users/david/work/website_watchdog/watchdog.py</string>
        <string>check</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/david/work/website_watchdog</string>
    
    <key>StandardOutPath</key>
    <string>/Users/david/work/website_watchdog/logs/service.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/david/work/website_watchdog/logs/service_error.log</string>
    
    <key>StartInterval</key>
    <integer>300</integer>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
